# kandan/kandan.conf
[program:kandan]
command=bundle exec thin start -e production
autostart=true
autorestart=true
username=kandan
directory=/srv/kandan
